# provisioning/roles/dm/templates/dm.service.j2
[Unit]
Description=DM
After=network.target

[Service]
User=dm
Group=www-data
WorkingDirectory=/home/dm
ExecStart=/usr/bin/java -jar /home/dm/{{ dm_jar }}

[Install]
WantedBy=multi-user.target
